# root/etc/services.d/tinyproxy/run
#!/command/with-contenv sh

echo "starting tinyproxy"
exec /usr/bin/tinyproxy -d; 

